# deploy/web/combsense-web.service.d/override.conf
[Service]
# Disable systemd sandboxing directives that fail 226/NAMESPACE in
# unprivileged LXCs. LXC itself provides the isolation.
PrivateTmp=no
PrivateDevices=no
ProtectSystem=no
ProtectHome=no
ProtectKernelTunables=no
ProtectKernelModules=no
ProtectControlGroups=no
RestrictAddressFamilies=
RestrictNamespaces=no
